# wafs_run/makeup_JVERF_GRID2GRID_WAFS.twind.cron
#!/bin/bash
set -xa
date

PDY1=20190723
PDY2=20190725

export DATA=/scratch4/NCEPDEV/stmp3/Yali.Mao/wafs.vrfy.com.twindonly
mkdir -p $DATA
cd $DATA

export vday=$PDY1
while [[ $vday < $PDY2 ]] ; do
  sh /scratch4/NCEPDEV/global/noscrub/Yali.Mao/git/save/wafs_run/makeup_JVERF_GRID2GRID_WAFS.twind.sh para para > ${vday}.log 2>&1
  vday=`$NDATE 24 ${vday}00 | cut -c1-8`
done

date
